# root/etc/s6-overlay/s6-rc.d/init-mod-universal-postgres-add-package/run
#!/usr/bin/with-contenv bash

# PostgreSQL 15 with VectorChord is pre-installed via docker mod
# This mod uses ghcr.io/tensorchord/vchord-postgres:pg15-v0.4.3
# All required system libraries are bundled in the mod
PGVERSION=15

# Ensure postgres user exists
if ! id -u postgres >/dev/null 2>&1; then
	echo "**** Creating postgres user ****"
	if [ -f /usr/bin/apt ]; then
		# Ubuntu
		useradd -r -s /bin/bash -d /var/lib/postgresql postgres
	elif [ -f /sbin/apk ]; then
		# Alpine
		adduser -D -H -s /bin/sh postgres
	fi
else
	echo "**** PostgreSQL user already exists ****"
fi

echo "**** PostgreSQL 15 with VectorChord binaries and dependencies included in mod ****"
