# root/etc/s6-overlay/s6-rc.d/init-mod-universal-postgres-setup/run
#!/usr/bin/with-contenv bash

# PostgreSQL 15 with VectorChord extension
PGVERSION=15

# for ubuntu
export PATH=$PATH:/usr/lib/postgresql/${PGVERSION}/bin/

# Check if trying to upgrade PostgreSQL version
if [ -f /config/postgres/PG_VERSION ] && [ $(cat /config/postgres/PG_VERSION) != ${PGVERSION} ]; then
	echo "Looks like you are trying to upgrade to PostgreSQL $PGVERSION, this will cause corruption and is not supported."
	echo "You may want to consider using the 'pg_upgrade' utility to upgrade your $(cat /config/postgres/PG_VERSION) to ${PGVERSION},"
	echo "otherwise stick with $(cat /config/postgres/PG_VERSION) by setting 'PG_VERSION' back to '$(cat /config/postgres/PG_VERSION)'"
	sleep infinity
fi

# Create postgres folders and set permissions
mkdir -p \
	/config/postgres/ \
	/run/postgresql/
chown -R postgres:postgres \
	/config/postgres/ \
	/run/postgresql/

# Initialize PostgreSQL database if not already initialized
if [ ! -f /config/postgres/PG_VERSION ]; then
	# Check if directory exists and is not empty (incomplete initialization)
	if [ -d /config/postgres ] && [ "$(ls -A /config/postgres)" ]; then
		echo "WARNING: /config/postgres exists but is not a valid database cluster"
		echo "Cleaning up incomplete initialization..."
		rm -rf /config/postgres/*
	fi

	echo "Initializing PostgreSQL Database with VectorChord extension"
	s6-setuidgid postgres initdb -D /config/postgres/ --username=postgres

	# Configure shared_preload_libraries for pg_stat_statements
	echo "shared_preload_libraries = 'pg_stat_statements'" >> /config/postgres/postgresql.conf

	s6-setuidgid postgres pg_ctl -D /config/postgres/ -o '-c listen_addresses= -p 5432' -w start
	s6-setuidgid postgres psql -c "ALTER USER postgres WITH PASSWORD 'postgres';"
	s6-setuidgid postgres psql -c "CREATE EXTENSION IF NOT EXISTS pg_stat_statements;"
	s6-setuidgid postgres psql -c "CREATE EXTENSION IF NOT EXISTS vectors;"
	s6-setuidgid postgres pg_ctl -D /config/postgres/ -m fast -w stop
else
	echo "PostgreSQL Database Already Initialised"
	# Ensure shared_preload_libraries is set correctly in existing installations
	if ! grep -q "shared_preload_libraries.*pg_stat_statements" /config/postgres/postgresql.conf; then
		echo "Updating shared_preload_libraries configuration"
		sed -i "/^shared_preload_libraries/d" /config/postgres/postgresql.conf
		echo "shared_preload_libraries = 'pg_stat_statements'" >> /config/postgres/postgresql.conf
	fi
fi
