# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-postgres/run
#!/usr/bin/with-contenv bash

# PostgreSQL 15 with VectorChord extension (Immich official image)
PGVERSION=15

# for ubuntu
export PATH=$PATH:/usr/lib/postgresql/${PGVERSION}/bin/

# Start PostgreSQL with shared_preload_libraries configuration
exec s6-setuidgid postgres \
	postgres -D /config/postgres/ \
	-c listen_addresses=localhost \
	-c port=5432 \
	-c shared_preload_libraries='pg_stat_statements,vchord,vectors.so'
